# root/etc/services.d/coredns/run
#!/usr/bin/with-contenv bash

UMASK_SET=${UMASK_SET:-022}

umask "$UMASK_SET"

exec /usr/bin/coredns -conf /app/Corefile